# familysafety-playwright/rootfs/etc/services.d/familysafety/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# S6 service finish script
# ==============================================================================

bashio::log.info "Family Safety Auth Service stopped"
